# iso-profile/slyme-os/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin slyme --noclear %I $TERM
Type=simple
