# iso-profile/slyme-os/airootfs/etc/systemd/system/slyme-ai.service
[Unit]
Description=Slyme OS AI Daemon
Documentation=https://slyme-os.github.io
After=ollama.service network.target
Wants=ollama.service

[Service]
Type=simple
User=root
Group=root
RuntimeDirectory=slyme-ai
RuntimeDirectoryMode=0755
ExecStart=/usr/local/bin/slyme-ai
Restart=on-failure
RestartSec=10
Environment=HOME=/root
Environment=SLYME_AI_SOCKET=/run/slyme-ai.sock
Environment=OLLAMA_HOST=127.0.0.1:11434
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slyme-ai

[Install]
WantedBy=multi-user.target
